# 7.4/apache/supervisor/conf.d/apache.conf
[program:http_server]
command=start-script
user=root
autorestart=true
autostart=true
